# supervisord-mysql.conf
[program:mysql]
command=/mysql-run.sh mysqld_safe
numprocs=1
autostart=true
autorestart=true
user=root
