# supervisord-tomcat.conf
[program:tomcat]
command=/tomcat-run.sh
numprocs=1
autostart=true
autorestart=false
user=tomcat
